# src/99-valve-index-reboot.rules
# Valve Index HMD: fix blank EDID after boot and suspend/resume.
# Works around kernel bug where the HMD appears as a 640x480 monitor.
# Sends a HID reboot command after hotplug or suspend/resume, at most once per boot per device.

KERNEL=="hidraw*", SUBSYSTEM=="hidraw", ATTRS{idVendor}=="28de", ATTRS{idProduct}=="2300", ACTION=="add|change", RUN+="/bin/sh -c '(sleep 2; test -f /run/.valve-index-rebooted-%k || { printf \\\\"\\\\\\\\x16\\\\\\\\x01\\\\" | cat - /dev/zero | head -c 64 > /dev/%k && touch /run/.valve-index-rebooted-%k; }) &'"
